# supervisord.conf
[supervisord]
nodaemon=true

[program:start.sh]
command=/usr/local/bin/start.sh
priority=1

[program:load.sh]
command=/load.sh
priority=2
